# cleanup-temp-files.service
[Unit]
Description=Cleanup temporary files for RecipeAI video-processor
After=network.target

[Service]
Type=oneshot
User=alan
WorkingDirectory=/home/alan/video-processor
ExecStart=/bin/bash /home/alan/video-processor/cleanup-temp-files.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
